# deploy/supervisor_iprs_backend.conf
[program:iprs-backend]
directory=/www/wwwroot/iprs/patent/backend
command=/www/wwwroot/iprs/patent/backend/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000
autostart=true
autorestart=true
startsecs=3
stopsignal=INT
redirect_stderr=true
stdout_logfile=/www/wwwroot/iprs/patent/backend/logs/uvicorn.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
user=www
environment=PYTHONUNBUFFERED="1"

